# cd-bus/fleet/systemd/cd-bus-subscriber@.service.d/onfailure.conf
# Drop-in for cd-bus-subscriber@.service: alert when the subscriber gives up
# (start-limit hit -> unit failed). A drop-in keeps the base unit untouched.
# %i propagates to the instance, so each service alerts for itself.
# Install at /etc/systemd/system/cd-bus-subscriber@.service.d/onfailure.conf
[Unit]
OnFailure=cd-bus-subscriber-alert@%i.service
