# cd-bus/fleet/systemd/cd-poll@.service.d/onfailure.conf
# Drop-in for cd-poll@.service: alert when the pull leg really fails. A drop-in
# keeps the base unit untouched. The base unit's SuccessExitStatus=0 1 means a
# single docker compose error is NOT failure — deploy.sh only exits 2 after 3
# consecutive failures, so this alert means the pull leg has been broken ~15 min.
# Install at /etc/systemd/system/cd-poll@.service.d/onfailure.conf
[Unit]
OnFailure=cd-poll-alert@%i.service
